# src/koreader-pkg/DEBIAN/postinst
#!/bin/sh

update-rc.d koreader defaults > /dev/null
