# files/rc.local
#!/bin/sh -e

depmod -a

DATA_PART=/dev/mmcblk1p3
if [ ! -e $DATA_PART ]; then
  DATA_PART=/dev/mmcblk0p3
fi
modprobe g_file_storage file=$DATA_PART

exit 0
